# packages/django-channel-tasks/django-channel-tasks-0.4.4.dev1.tar.gz/django-channel-tasks-0.4.4.dev1/nginx-unit/channel-tasks.template.service
[Unit]
Description=channel-tasks

[Service]
Type=simple
Restart=always
RestartSec=5s
ExecStart=${CHANNEL_TASKS_BIN_DIR}/restart-channel-tasks.sh
Requires=docker.service
Environment="DOCKER_ARCH_PREFIX=${DOCKER_ARCH_PREFIX}"
Environment="CHANNEL_TASKS_DOCKER_HOME=${CHANNEL_TASKS_DOCKER_HOME}"
Environment="CHANNEL_TASKS_PACKAGE=${CHANNEL_TASKS_PACKAGE}"
Environment="CHANNEL_TASKS_LOG_LEVEL=${CHANNEL_TASKS_LOG_LEVEL}"
Environment="CHANNEL_TASKS_BIN_DIR=${CHANNEL_TASKS_BIN_DIR}"
Environment="CHANNEL_TASKS_INI_PATH=${CHANNEL_TASKS_INI_PATH}"
Environment="TASK_ADMIN_USER=${TASK_ADMIN_USER}"
Environment="TASK_ADMIN_EMAIL=${TASK_ADMIN_EMAIL}"

[Install]
WantedBy=multi-user.target
